# setup/gs.conf
description "gs server"
author  "jt"

start on (local-filesystems and net-device-up IFACE=eth0)
stop on shutdown

respawn

exec sudo -u node sh -c "/usr/local/bin/node /var/local/sites/gs/app.js >> /var/log/gs.log 2>&1"

